# keycloak/rootfs/etc/services.d/keycloak/run
#!/usr/bin/env bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting keycloak..."
exec "/var/app/bin/kc.sh" start --optimized
